# deploy/systemd/ai-internships.timer
[Unit]
Description=Run AI Internship Tracker hourly

[Timer]
OnCalendar=hourly
Persistent=true
Unit=ai-internships.service

[Install]
WantedBy=timers.target
